# ansible_collections/ska_collections/monitoring/roles/snmp/files/snmpd.service
[Service]

ExecStart=

ExecStart=/usr/sbin/snmptrapd -Ln -f -Lf /var/log/snmptrapd.log