# ansible_collections/ska_collections/monitoring/roles/snmp/templates/snmpd.service.j2
[Service]

ExecStart=

ExecStart=/usr/sbin/snmptrapd -Ln -f -Lf {{ snmp_trap_logs_dir }}/snmptrapd.log